# scripts/bg_tasks.cron
# 计划任务参考样本
# 备注
# 1.不使用sudo
# 2.需要设定环境路径变量
# 3.运行脚本（使用激活环境后再用命令行，无法正确执行）
# 4."0"代表星期天

# 运行方式 $crontab -e <path><to><crontab.cron>

# m h  dom mon dow   command

# 在所在的python环境运行 echo $PATH，拷贝设置PATH环境变量
PATH=/home/ldf/anaconda3/envs/fsweb/bin:/home/ldf/anaconda3/bin:/home/ldf/bin:/home/ldf/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin

# 24小时财经新闻（每分钟执行一次）
* * * * * python /home/ldf/cswd/scripts/per_minutely_tasks.py >> /home/ldf/stockdata/logs/global_news.log 2>&1

# 融资融券（每周一到周六早上9：15点执行一次）
15 09 * * 1-6 python /home/ldf/cswd/scripts/am_9_margindata.py >> /home/ldf/stockdata/logs/margin.log 2>&1

# 财务公告，权益变动（每天8、12、16、20点执行一次）
0 08,12,16,20 * * * python /home/ldf/cswd/scripts/hourly_tasks.py >> /home/ldf/stockdata/logs/info.log 2>&1

# 分时交易（每周一到五交易时段内每分钟执行一次）
25-59 09 * * 1-5 python /home/ldf/cswd/scripts/minutely_tasks.py >> /home/ldf/stockdata/logs/live_trading.log 2>&1
*     10 * * 1-5 python /home/ldf/cswd/scripts/minutely_tasks.py >> /home/ldf/stockdata/logs/live_trading.log 2>&1
0-30  11 * * 1-5 python /home/ldf/cswd/scripts/minutely_tasks.py >> /home/ldf/stockdata/logs/live_trading.log 2>&1
*  13-14 * * 1-5 python /home/ldf/cswd/scripts/minutely_tasks.py >> /home/ldf/stockdata/logs/live_trading.log 2>&1

# 实时报价（交易时段内提取实时报价数据，在内部使用循环，每分钟运行60次）
15-59    9 * * 1-5 python /home/ldf/cswd/scripts/live_quote.py >> /home/ldf/stockdata/logs/quote.log 2>&1 
*       10 * * 1-5 python /home/ldf/cswd/scripts/live_quote.py >> /home/ldf/stockdata/logs/quote.log 2>&1
0-30    11 * * 1-5 python /home/ldf/cswd/scripts/live_quote.py >> /home/ldf/stockdata/logs/quote.log 2>&1
*    13-14 * * 1-5 python /home/ldf/cswd/scripts/live_quote.py >> /home/ldf/stockdata/logs/quote.log 2>&1

# 盘后日线数据（每周一到五18点开始刷新）
0 18 * * 1-5 python /home/ldf/cswd/scripts/daily_tasks.py >> /home/ldf/stockdata/logs/daily_trading_data.log 2>&1
